# .docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:mails]
directory=/var/www/html
command=php bin/console messenger:consume mail --time-limit=3600
autostart=true
autorestart=true

[program:async]
directory=/var/www/html
command=php bin/console messenger:consume async --time-limit=3600 --memory-limit=256M
autostart=true
autorestart=true

[program:checks]
directory=/var/www/html
command=php bin/console messenger:consume checks --time-limit=3600 --memory-limit=256M
autostart=true
autorestart=true