# .docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:mails]
directory=/var/www/html
command=php bin/console messenger:consume mail --time-limit=3600
autostart=true
autorestart=true

[program:async]
directory=/var/www/html
command=php bin/console messenger:consume async --time-limit=3600 --memory-limit=256M
autostart=true
autorestart=true

[program:checks]
directory=/var/www/html
command=php bin/console messenger:consume scheduler_default --time-limit=3600 --memory-limit=256M
autostart=true
autorestart=true

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:cron]
command=crond -f
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0